# .config/systemd/user/gBar.service
[Unit]
Description=gBar

[Service]
ExecStart=/usr/bin/gBar bar 0
Restart=always

[Install]
WantedBy=default.target
